# system/nilpferd.target
[Unit]
Description=Nilpferd - self-hosted cloud by malyvsen.
Requires=nilpferd-caddy.service
Requires=nilpferd-filebrowser.service
